# alibaba_crawler.ini
[program:alibaba_crawler_1]
environment = LD_LIBRARY_PATH="/usr/local/lib:$LD_LIBRARY_PATH", PYTHONPATH="$PYTHONPATH:/usr/local/lib/python3.7/site-packages/alibaba_crawler"
command = /usr/local/bin/python3.7 /usr/local/bin/run_alibaba_crawler 
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /logs/alibaba_crawler/alibaba_crawler_1.log

[group:alibaba_crawler]
programs=alibaba_crawler_1
